# files/supervisor.webvirtmgr.conf
[program:webvirtmgr]
command=/usr/bin/python /opt/webvirtmgr/manage.py run_gunicorn -c /opt/webvirtmgr/conf/gunicorn.conf.py
directory=/opt/webvirtmgr
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/webvirtmgr.log
redirect_stderr=true
user=www-data

[program:webvirtmgr-console]
command=/usr/bin/python /opt/webvirtmgr/console/webvirtmgr-console
directory=/opt/webvirtmgr
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/webvirtmgr-console.log
redirect_stderr=true
user=www-data
